# Projects/private-ca/roles/compute/files/nomad.service
[Unit]
Description=Nomad

[Service]
Type=simple
ExecStart=/usr/bin/nomad agent -config=/etc/nomad

[Install]
WantedBy=multi-user.target
